# tools/ci/docker/tessia-server/assets/supervisord.conf
[supervisord]
nodaemon=true
logfile_maxbytes=20MB
logfile_backups=3

[program:tessia-api]
command=uwsgi --yaml /etc/tessia/uwsgi.yaml
redirect_stderr=true
stdout_logfile=/var/log/tessia/uwsgi.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3

[program:tessia-scheduler]
command=tess-scheduler
redirect_stderr=true
stdout_logfile=/var/log/tessia/scheduler.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
user=tessia
environment=HOME="/home/tessia",USER="tessia"

[program:installer-webhook]
command=installer-webhook
redirect_stderr=true
stdout_logfile=/var/log/tessia/installer-webhook.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stopsignal=TERM
stopwaitsecs=3
user=tessia
environment=HOME="/home/tessia",USER="tessia"
